# conf/supervisord/laravel-listener.conf
[program:laravel-listener]
command=php artisan queue:listen
directory=/var/www/laravel
stdout_logfile=/var/log/laravel-listener.log
redirect_stderr=true